# supervisor.conf
[supervisord]
nodaemon=true
user=root

[program:web]
command=/bin/sh -c 'if [ -n "$PORT" ]; then exec php -S 0.0.0.0:$PORT -t /var/www/html/web /var/www/html/web/index.php; else exec php-fpm; fi'
autostart=true
autorestart=true

[program:cron]
command=cron -f
autostart=true
autorestart=true
