# supervisor.conf
[supervisord]
nodaemon=true
user=root

[program:web]
command=/bin/sh -c 'if [ "${RUN_PHP_FPM:-false}" = "true" ]; then exec php-fpm; else exec php -S 0.0.0.0:${PORT:-8080} -t /var/www/html/web /var/www/html/web/router.php; fi'
autostart=true
autorestart=true

[program:cron]
command=cron -f
autostart=true
autorestart=true
